# contrib/init/bitcoind.conf
description "Actinium Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env ACTINIUMD_BIN="/usr/bin/Actiniumd"
env ACTINIUMD_USER="actinium"
env ACTINIUMD_GROUP="actinium"
env ACTINIUMD_PIDDIR="/var/run/Actiniumd"
# upstart can't handle variables constructed with other variables
env ACTINIUMD_PIDFILE="/var/run/Actiniumd/Actiniumd.pid"
env ACTINIUMD_CONFIGFILE="/etc/actinium/Actinium.conf"
env ACTINIUMD_DATADIR="/var/lib/Actiniumd"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # ActiniumD will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$ACTONIUMD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run Actiniumd."
        echo "The setting must appear in $ACTINIUMD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Actinium Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$ACTINIUMD_PIDDIR"
    chmod 0755 "$ACTINIUMD_PIDDIR"
    chown $ACTINIUMD_USER:$ACTINIUMD_GROUP "$ACTINIUMD_PIDDIR"
    chown $ACTINIUMD_USER:$ACTINIUMD_GROUP "$ACTINIUMD_CONFIGFILE"
    chmod 0660 "$ACTINIUMD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$ACTINIUMD_PIDFILE" \
    --chuid $ACTINIUMD_USER:$ACTINIUMD_GROUP \
    --exec "$ACTINIUMD_BIN" \
    -- \
    -pid="$ACTINIUMD_PIDFILE" \
    -conf="$ACTINIUMD_CONFIGFILE" \
    -datadir="$ACTINIUMD_DATADIR" \
    -disablewallet \
    -daemon

